# systemd/user/rencrow-data-weekly.service
[Unit]
Description=RenCrow data weekly feature snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PICOCLAW_REPO_DIR=@PICOCLAW_REPO_DIR@
ExecStart=/usr/bin/env bash -lc 'exec "%h/.local/share/picoclaw/scripts/rencrow_data_scheduler.sh" weekly'
WorkingDirectory=%h
StandardOutput=append:%h/.picoclaw/logs/rencrow-data-weekly.log
StandardError=append:%h/.picoclaw/logs/rencrow-data-weekly.log

[Install]
WantedBy=default.target
