# systemd/user/rencrow-resilience.timer
[Unit]
Description=Check RenCrow CORE liveness and incident recovery

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true
Unit=rencrow-resilience.service

[Install]
WantedBy=timers.target
