# systemd/user/rencrow.service.d/10-panic-stack.conf
[Service]
# panic時にpanicしたgoroutineだけでなく、全goroutineのstackをjournalへ出力する。
Environment=GOTRACEBACK=all
StandardOutput=journal
StandardError=journal
# 高頻度ログ中でもpanicやerrorをjournal側で間引かない。
LogRateLimitIntervalSec=0
LogRateLimitBurst=0
