# systemd/user/rencrow.service.d/20-resilience.conf
[Unit]
# RenCrowは連続異常終了時にもsystemdのstart limitで停止しない。
StartLimitIntervalSec=0

[Service]
# 正常停止は無視し、panic・signal・timeoutなどの終了証拠だけを事故台帳へ保存する。
ExecStopPost=-%h/.local/bin/rencrow resilience capture-stop
