# systemd/user/rencrow.service.d/30-games-observer.conf
# RenCrow_GAMES 共有 observer (ランチャー) の接続先。
# /viewer/games/launch と observer proxy の転送先に使われる
# (docs/06_Public_API仕様.md Game Launch 節、マルチペルソナ WP5 A-2)。
# observer を同一ホストで常駐させる構成が既定。別ホスト (例: Windows 機) で
# 動かす場合は URL をそのホストへ差し替え、observer 側の -observer-addr を
# 0.0.0.0 バインドにすること。
[Service]
Environment=RENCROW_GAMES_OBSERVER_URL=http://127.0.0.1:18796
